# services/lizard-brain-www/run
#!/bin/nosh
chdir ../..
bin/config-set-env
tcp-socket-listen 127.0.0.1 6000
tcp-socket-accept --no-delay
cgid plackup www/cgi-bin/impulse-www
