# services/lizard-brain-www/run
#!/bin/dash
cd ../..
plackup bin/impulse-www -p 6000
