# services/lizard-brain-www/run
#!/bin/nosh
chdir ../..
bin/config-set-env
plackup -l 127.0.0.1:6000 ./bin/httpd
