# services/lizard-brain-twitter/run
#!/bin/nosh
chdir ../..
envdir --ignore-nodir --chomp /home/frew/.lizard-brain
bin/config-set-env
bin/impulse-twitter
