# services/lizard-brain/run
#!/bin/nosh
chdir ../..
bin/config-set-env
bin/service
